# crontab
*/4 * * * * curl -v  http://localhost:5000/api/elecciones/download_data /home/stanlee321/Desktop/2020/elecciones/elecciones_scraper_2020/app.log 2>&1
